# legacy/admin/ka9q-admin.service
[Unit]
Description=KA9Q Web SDR Admin Dashboard
After=network-online.target
Wants=ka9q-web.service

[Service]
Type=simple
User=radio
Group=radio
ExecStart=/home/wsprdaemon/wsprdaemon/ka9q-web/admin/venv/bin/python admin.py
WorkingDirectory=/home/wsprdaemon/wsprdaemon/ka9q-web/admin
Environment=KA9Q_ADMIN_CONF=/home/wsprdaemon/wsprdaemon/ka9q-web/admin/admin.conf
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
